# DEBIAN/postinst
#!/bin/bash
update-rc.d run-once defaults
